# cron/crontab
SHELL=/bin/bash
HOME=/home/ubuntu
* * * * * /usr/bin/runrescuetime 2>&1 >> /tmp/rescuetime.log
